# profile.d/perso.sh
# ==================================
#     ALIAS
# ==================================
alias runepicurom='~/bin/google_appengine/dev_appserver.py ~/sandboxes/epicurom/'
alias epicurom='cd ~/sandboxes/epicurom'
alias mytr='vi ~/Dropbox/tri/training/week/this_week.txt'
alias mytr_dir='cd ~/Dropbox/tri/training/week'

