# profile.d/perso.sh
# ==================================
#     ALIAS
# ==================================
alias runepicurom='~/bin/google_appengine/dev_appserver.py ~/sandboxes/epicurom/'
alias epicurom='cd ~/sandboxes/epicurom'

